# roles/tqm/templates/tqm.timer.j2
[Unit]
Description=tqm timer

[Timer]
OnUnitInactiveSec=60min


[Install]
WantedBy=timers.target

